# config/.config/systemd/user/change_wallpaper.timer
[Unit]
Description=Timer para alterar papéis de parede a cada hora

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target